# deploy/the-controller.service
# User-level systemd unit for The Controller server mode.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp deploy/the-controller.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now the-controller
#
# The service runs as your user — no sudo, no permission conflicts.
# Requires lingering so it survives logout:
#   loginctl enable-linger $USER
#
# ExecStart uses the user's login shell so the server inherits the full
# environment (.profile/.zshrc/.zprofile) — this is how it finds claude,
# tmux, nvm-managed node, cargo-installed tools, etc.
# The deploy script replaces LOGIN_SHELL with the actual shell path.

[Unit]
Description=The Controller - Claude Code Session Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.the-controller
EnvironmentFile=%h/.the-controller/server.env
ExecStart=LOGIN_SHELL -lc 'exec %h/.the-controller/server'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
